# packages/systemd/systemd.install
post_install() {   
systemd-machine-id-setup &> /dev/null
systemctl preset-all &> /dev/null
}

